# midi-main.service
sudo tee /etc/systemd/system/midi-main.service >/dev/null <<'UNIT'
[Unit]
Description=MIDI Tracker Main App (GUI on :0)
# Wait for the desktop to be up so Tk has a display
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=root
Group=swag
WorkingDirectory=/home/swag/Documents/git/midiTrackerApp

# Use project venv
Environment="VIRTUAL_ENV=/home/swag/Documents/git/midiTrackerApp/.venv"
Environment="PATH=/home/swag/Documents/git/midiTrackerApp/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Environment=PYTHONUNBUFFERED=1

# Point to user's desktop session
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/swag/.Xauthority
# Helpful for some libs on Wayland/X11 bridges:
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus

ExecStart=/home/swag/Documents/git/midiTrackerApp/.venv/bin/python /home/swag/Documents/git/midiTrackerApp/main.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical.target
UNIT

sudo systemctl daemon-reload
sudo systemctl enable --now midi-main.service
systemctl status midi-main.service --no-pager
